# com.acsandmann.swipe.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.acsandmann.swipe</string>

    <key>ProgramArguments</key>
    <array>
      <string>@TARGET_PATH@</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>

	<key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/tmp/swipe.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/swipe.err</string>
  </dict>
</plist>
